# debian/do-backup.postrm
#!/bin/sh
set -e

invoke-rc.d do-backup stop > /dev/null
update-rc.d -f do-backup remove > /dev/null

#DEBHELPER#
