# debian/do-backup.postinst
#!/bin/sh
set -e

update-rc.d do-backup defaults > /dev/null
invoke-rc.d do-backup start > /dev/null

#DEBHELPER#
